# configs/speedy-supervisor.ini
[program:speedy-api]
command=python3.5 api.py
directory=/Users/fahim/Source/speedy
autostart=true
autorestart=true
stdout_logfile=/usr/local/var/log/speedy-api.log
redirect_stderr=true

[program:speedy-scheduler]
command=python3.5 app.py
directory=/Users/fahim/Source/speedy
autostart=true
autorestart=true
stdout_logfile=/usr/local/var/log/speedy.log
redirect_stderr=true
